# watchmail.service
[Unit]
Description=Watch for emails

[Service]
Type=simple
ExecStart=
User=

[Install]
WantedBy=multi-user.target
